# systemConf/pacupdate.service
[Unit]
Description=Fetch Packages to Update

[Service]
User=root
ExecStart=/usr/bin/pacman -Sy > /dev/null

